# quotsi/sh/runqotd.service
[Unit]
Description=Run the daily qotd

[Service]
ExecStart=/home/watsnav/scripts/www/watsnav.github.io/quotsi/sh/runqotd.sh

