# quotsi/sh/runqotd-service/runqotd.timer
[Unit]
Description=Runs runqotd daily

[Timer]
OnUnitActiveSec=1d
Unit=runqotd.service

[Install]
WantedBy=multi-user.target
